# arch/system/greetd/update-issue.service
[Unit]
Description=Refresh terminal login system snapshot
Before=greetd.service

[Service]
Type=oneshot
ExecStart=/usr/local/lib/workstation-identity/update-issue

[Install]
WantedBy=greetd.service
